# production_scripts/etc/systemd/system/forenings_medlemmer_backup.timer
[Unit]
Description=Run django backup timer 

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hour

[Install]
WantedBy=basic.target
